# roles/darkedges.multichain/templates/etc/profile.d/multichain.sh.j2
export PATH=$PATH:{{ multichain_usr_dir }}/:$PATH:{{ multichain_usr_dir }}/bin